# systemd/pidial.service
[Unit]
Description=PiDial Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=prcutler
ExecStart=/home/prcutler/workspace/pi-dial/.venv/bin/python /home/prcutler/workspace/pi-dial/pidial.py
ExecStartPost=
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target